# crontab
0 3 * * * /bin/bash /app/script.sh >> /var/log/cron.log 2>&1